# Bash/gunicorn.socket
[Unit]
Description=Gunicorn.sock Live...!

[Socket]
ListenStream=/run/gunicorn.sock

[Install]
WantedBy=sockets.target
